# deploy/fnk0043.service
[Unit]
Description=FNK0043 Smart Car Web Controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/robotic-projects
Environment=HOME=%h
Environment=PATH=%h/robotic-projects/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/robotic-projects/.venv/bin/fnk0043-server
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
