# deploy/fnk0043.service
[Unit]
Description=FNK0043 Smart Car Web Controller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/robotic-projects
Environment=HOME=%h
Environment=FNK0043_PARAMS=%h/robotic-projects/params.json
Environment=PATH=%h/robotic-projects/.venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=%h/robotic-projects/.venv/bin/fnk0043-server
Restart=on-failure
RestartSec=5
TimeoutStopSec=15
KillMode=mixed

[Install]
WantedBy=default.target
